# etc/crontab
0 5  * * *  /root/XRAY-RELEASE/cron_job.sh
